# automation/systemd/textjepa-research-watch.timer
[Unit]
Description=Run the TextJEPA research watcher every two minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target

